# images/frr.ini
[program:frr]
directory=/
command=/usr/lib/frr/frrinit.sh start
user=root
autostart=true
autorestart=false
stdout_logfile=/frr_log
stderr_logfile=/frr_log
redirect_stdout=true
redirect_stderr=true
