# images/lldpd.ini
[program:lldpd]
directory=/
command=/usr/sbin/lldpd -d -I !eth0,*
user=root
autostart=true
autorestart=true